# build_dir/target-mips_24kc_musl/root-ar71xx/etc/uci-defaults/luci-i18n-ntpc-zh-cn
uci set luci.languages.zh_cn='中文 (Chinese)'; uci commit luci
