# build_dir/target-mips_24kc_musl/root.orig-ar71xx/usr/lib/opkg/info/monit.list
/etc/init.d/monit
/usr/bin/monit
/etc/monitrc
