# build_dir/target-mips_24kc_musl/root.orig-ar71xx/usr/lib/opkg/info/nlbwmon.list
/usr/sbin/nlbw
/etc/config/nlbwmon
/usr/share/nlbwmon/protocols
/etc/init.d/nlbwmon
/usr/sbin/nlbwmon
