# build_dir/target-mips_24kc_musl/root-ar71xx/usr/lib/opkg/info/mcproxy.list
/etc/mcproxy.conf
/etc/config/mcproxy
/etc/init.d/mcproxy
/usr/sbin/mcproxy
