# build_dir/target-mips_24kc_musl/root-ar71xx/usr/lib/opkg/info/wifidog.list
/usr/lib/libhttpd.so.0.0.0
/usr/bin/wifidog
/etc/init.d/wifidog
/usr/lib/libhttpd.so.0
/usr/bin/wdctl
/etc/wifidog-msg.html
/etc/wifidog.conf
